# .config/kanata-tray/kanata-tray.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST File Format//EN" "https://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rszyma.kanata-tray</string>
    <key>ProgramArguments</key>
    <array>
    <string>/Users/govind.r/.local/bin/kanata-tray</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Library/Frameworks/GStreamer.framework/Versions/1.0/bin:/opt/homebrew/opt/postgresql@12/bin:/opt/homebrew/opt/coreutils/libexec/gnubin:/opt/homebrew/bin:/opt/homebrew/sbin:/Users/govind.r/.pyenv/bin:/Users/govind.r/.pyenv/shims:/Users/govind.r/.pyenv/bin:/usr/local/bin:/System/Cryptexes/App/usr/bin:/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/MacGPG2/bin:/Library/Apple/usr/bin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/local/bin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/bin:/var/run/com.apple.security.cryptexd/codex.system/bootstrap/usr/appleinternal/bin:/opt/homebrew/opt/postgresql@12/bin:/opt/homebrew/opt/coreutils/libexec/gnubin:/opt/homebrew/bin:/opt/homebrew/sbin:/Users/govind.r/.nix-profile/bin:/nix/var/nix/profiles/default/bin:/Users/govind.r/.pyenv/bin:/Applications/kitty.app/Contents/MacOS:/Users/govind.r/.local/bin/statusbar:/Users/govind.r/.local/bin/private/dev-scripts:/Users/govind.r/.local/bin/private:/Users/govind.r/.local/bin/tmux:/Users/govind.r/.local/bin/cron:/Users/govind.r/.local/bin:/Users/govind.r/.local/share/npm/bin:/Users/govind.r/.local/share/go/bin:/Users/govind.r/Library/Android/sdk/emulator:/Users/govind.r/Library/Android/sdk/cmdline-tools/latest/bin:/Users/govind.r/Library/Android/sdk/platform-tools:/Users/govind.r/.cargo/bin:/Users/govind.r/.local/bin/statusbar:/Users/govind.r/.local/bin/private/dev-scripts:/Users/govind.r/.local/bin/private:/Users/govind.r/.local/bin/tmux:/Users/govind.r/.local/bin/cron:/Users/govind.r/.local/bin:/Users/govind.r/.local/share/npm/bin:/Users/govind.r/.local/share/go/bin:/Users/govind.r/Library/Android/sdk/emulator:/Users/govind.r/Library/Android/sdk/cmdline-tools/latest/bin:/Users/govind.r/Library/Android/sdk/platform-tools:/Users/govind.r/.cargo/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
 	     <false/>
 	     <key>Crashed</key>
 	     <true/>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/kanata-tray.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/kanata-tray.err.log</string>
    <key>ProcessType</key>
    <string>Interactive</string>
    <key>Nice</key>
    <integer>-20</integer>
</dict>
</plist>
